# systemd/notion-mcp-remote.service
[Unit]
Description=Notion MCP Remote Server
After=network.target

[Service]
Type=simple
User=ldraney
WorkingDirectory=/home/ldraney/notion-mcp-remote
EnvironmentFile=/home/ldraney/notion-mcp-remote/.env
ExecStart=/home/ldraney/notion-mcp-remote/.venv/bin/python server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
